# gtask_admin/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:uwsgi]
command=/bin/bash -c "uwsgi --ini ni_uwsgi.ini:docker"