# supervisord.conf
[supervisord]
nodaemon=true

[program:uwsgi]
command=/bin/bash -c "uwsgi --ini ni_uwsgi.ini:docker"
autorestart=false
redirect_stderr=true